# etc/init.d/vnc
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	need net
}



start() {
	ebegin "Starting VNC server"
	/bin/tcpulse /etc/cert.pem 8081 > /dev/null 2>&1
	eend $?
}

stop() {
	ebegin "Stopping VNC server"
        pkill -U 1000

	rm -rf /tmp/.X1-lock  > /dev/null 2>&1
        rm -rf /tmp/.X11-unix > /dev/null 2>&1

 	eend 0
}

restart() {
        svc_stop
        svc_start
}
