# etc/init.d/vnc
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	need net
}



start() {
	ebegin "Starting VNC server"
	startx
	/bin/tcpulse /etc/cert.pem 5701
	source /etc/profile && export UNIXPW_DISABLE_SSL=1 && export UNIXPW_DISABLE_LOCALHOST=1 && cd ~ && export FD_GEOM=1920x1080 && x11vnc -shared -many -nomodtweak -noxkb -afteraccept "/bin/update_allowed"  -rfbport 5801 -unixpw_cmd /bin/verify_pw -loop -display ":1" -ssl /etc/cert.pem -repeat -clear_all -skip_lockkeys
	eend $?
}

stop() {
	ebegin "Stopping VNC server"
        pkill -U 1000

	rm -rf /tmp/.X1-lock 2>/dev/null
        rm -rf /tmp/.X11-unix >/dev/null 2>&1

 	eend 0
}

restart() {
        svc_stop
        svc_start
}
